# supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as in [*_http_server] if set
;password=123                ; should be same as in [*_http_server] if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available


;Commands to run the projects are made for respective programs
;Please update the file directory accordingly

[program:jupyter]
command=jupyter notebook --ip='*'
stdout_logfile=/var/log/lab/jupyter.log

[program:video-spooler]
command=python3 /home/bdt1/projects/videoprocessor/bdt/kafka/video/processor/__init__.py
environment=PYTHONPATH="/home/bdt1/projects/videoprocessor"
stdout_logfile=/var/log/lab/video-spooler.log

[program:video-app]
command=python3 /home/bdt1/projects/videoprocessor/bdt/app.py
environment=PYTHONPATH="/home/bdt1/projects/videoprocessor"
stdout_logfile=/var/log/lab/video-app.log

[program:plainframe-streamer]
command=spark-submit --master yarn plainframesaverlg.jar       ;This is the plain frame streamer jar file exported into a single binary.
stdout_logfile=/var/log/lab/plainframesaver.log

[program:rgbsplitter-streamer]
command=spark-submit --jars /opt/spark-2.4.4/jars/spark-streaming-kafka-0-8-assembly_2.11-2.4.0.jar --master  yarn  rgbstreamer/RDBSplitter.py
stdout_logfile=/var/log/lab/rgbsplitter.log

[program:yolo-streamer]
command=spark-submit --jars /opt/spark-2.4.4/jars/spark-streaming-kafka-0-8-assembly_2.11-2.4.4-sources.jar --master yarn videoprocessor/yolostream/yolo_stream.py localhost:2181 video-stream
stdout_logfile=/var/log/lab/yolo-streamer.log

